# vps-backup/systemd/rclone-backup.service
[Unit]
Description=Rclone-backup
After=network.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/bin/bash /root/run-rclone-backup.sh
