# auto-wlt/wlt.service
[Unit]
Description=Auto Wlt
After=network.target network-online.target

[Service]
Type=oneshot
ExecStart=/home/pi/wlt.sh
