# DEPRECATED/auto-wlt/wlt.service
[Unit]
Description=Auto WLT
After=network.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/bin/bash /home/pi/wlt.sh
